# apaas-locust-rtc-demo/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
logfile_maxbytes=10MB
logfile_backups=3
loglevel=info

[program:go-demo-server]
command=/opt/go_deploy_demo/agora_sdk_server_linux
directory=/opt/go_deploy_demo
autostart=true
autorestart=true
startsecs=3
startretries=3
stdout_logfile=/var/log/go-demo-server.log
stdout_logfile_maxbytes=10MB
stderr_logfile=/var/log/go-demo-server.err
stderr_logfile_maxbytes=5MB
priority=10

[program:locust]
command=/opt/entrypoint.sh
directory=/opt/apaas_scene
autostart=false
autorestart=false
startsecs=5
startretries=1
stdout_logfile=/var/log/locust.log
stdout_logfile_maxbytes=50MB
stderr_logfile=/var/log/locust.err
stderr_logfile_maxbytes=10MB
priority=20
